# scripts/com.zuby.nyc-events.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.zuby.nyc-events</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-lc</string>
        <string>cd "/Users/zuby/AI Building/nyc-events-agent" &amp;&amp; ./.venv/bin/python main.py --only eventbrite_kids</string>
    </array>

    <!-- Run every night at 2:00 AM. -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>2</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/zuby/AI Building/nyc-events-agent/logs/launchd.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/zuby/AI Building/nyc-events-agent/logs/launchd.err.log</string>
</dict>
</plist>
